# mqtt_subscriber/files/mqtt_subscriber.init
#!/bin/sh /etc/rc.common

USE_PROCD=1

start_service() {

}

reload_service() {
    stop
    start
}